# config/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python /app/app.py
autostart=true
autorestart=true

[program:listener]
command=python /app/listener.py
autostart=true
autorestart=true
